# example/prom-elastic-alert-supervisor.ini
[program:prom-elastic-alert]
command=/opt/prom-elastic-alert/prom-elastic-alert --config /opt/prom-elastic-alert/config.yaml
directory=/opt/prom-elastic-alert
stdout_logfile=/var/log/prom-elastic-alert.log
autostart=true
autorestart=true
redirect_stderr=true
stopsignal=QUIT
user=prometheus
